# ptt.service
# ============================================================
# PTT Advertising — Gunicorn (systemd)
#
# Cách dùng trên VPS (Ubuntu):
#   sudo cp ptt.service /etc/systemd/system/ptt.service
#   sudo nano /etc/systemd/system/ptt.service   # sửa User, WorkingDirectory, cổng -b
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ptt
#   sudo systemctl status ptt
#
# Nginx site 2 phải proxy_pass tới cùng cổng với -b (vd. 127.0.0.1:8002).
# Trong thư mục app cần file .env với ít nhất: FLASK_SECRET_KEY=...
# ============================================================

[Unit]
Description=PTT Advertising — Flask (Gunicorn)
After=network.target

[Service]
Type=simple
User=deploy
Group=www-data
WorkingDirectory=/var/www/ptt
Environment=PATH=/var/www/ptt/.venv/bin
EnvironmentFile=-/var/www/ptt/.env
ExecStart=/var/www/ptt/.venv/bin/gunicorn -c gunicorn.conf.py -w 3 -b 127.0.0.1:8002 --timeout 120 app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
