# deploy/ptt-meta-clickhouse-export.timer
[Unit]
Description=PTT — timer Meta/Google daily facts → ClickHouse (B14)
Requires=ptt-meta-clickhouse-export.service

[Timer]
OnCalendar=*-*-* 04:45:00 Asia/Ho_Chi_Minh
Persistent=true

[Install]
WantedBy=timers.target
