# deploy/ptt-ops-web.service
[Unit]
Description=PTT Ops — Next.js standalone (internal staff)
After=network.target ptt-crm-api.service

[Service]
Type=simple
User=deploy
Group=www-data
WorkingDirectory=/var/www/ptt/services/ops-web
Environment=NODE_ENV=production
Environment=PORT=3200
EnvironmentFile=-/var/www/ptt/.env
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
