# deploy/ptt-wave1-leads-soak.timer
[Unit]
Description=Daily Wave 1 CRM leads soak snapshot

[Timer]
OnCalendar=*-*-* 06:15:00
Persistent=true

[Install]
WantedBy=timers.target
