# deploy/ptt-seo-ga4-sync.timer
[Unit]
Description=PTT — timer SEO GA4 OAuth sync hàng ngày 03:30 ICT (Phase 4)

[Timer]
OnBootSec=25min
OnCalendar=*-*-* 03:30:00
Timezone=Asia/Ho_Chi_Minh
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
