# deploy/ptt-email-journey.timer
[Unit]
Description=PTT — timer email journey execution every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
