# deploy/ptt-lead-review-queue-sync.timer
[Unit]
Description=Hourly B2 review queue sync (Wave B4)

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
