# deploy/ptt-seo-gsc-sync.timer
[Unit]
Description=PTT — timer SEO GSC OAuth sync hàng ngày 03:00 ICT (Phase 4)

[Timer]
OnBootSec=20min
OnCalendar=*-*-* 03:00:00
Timezone=Asia/Ho_Chi_Minh
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
